# scrapegoat.service
[Unit]
Description=Scrapegoat Documentation Scraping Service
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
Documentation=http://git.den.lan/pub/scrapegoat

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/scrapegoat
ExecStart=/usr/bin/docker compose -f docker-compose.byo-postgres.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.byo-postgres.yml down
ExecReload=/usr/bin/docker compose -f docker-compose.byo-postgres.yml restart
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300
TimeoutStopSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scrapegoat

[Install]
WantedBy=multi-user.target
